# ansible/roles/k3s/templates/k3s.slice.j2
[Unit]
Description=Dedicated slice for K3s workloads

[Slice]
CPUAccounting=true
MemoryAccounting=true
